# udev-rules/90-redflat-devicekit-usb.rules
# Make sure to change the `fx` to you user name.
ACTION=="add", SUBSYSTEM=="usb", RUN+="/usr/bin/devicekit_runner.sh fx usbAdded string:$ID_VENDOR"
